# deploy/auto-translate-bot.service
[Unit]
Description=Auto-Translade Telegram Bot
After=network-online.target

[Service]
Type=simple
User=hai
WorkingDirectory=/home/hai/Auto-Translade-video
EnvironmentFile=/home/hai/Auto-Translade-video/.env
ExecStart=/home/hai/Auto-Translade-video/.venv/bin/python -m src.telegram_bot
Restart=on-failure
RestartSec=5s
StandardOutput=append:/var/log/auto-translate-bot.log
StandardError=append:/var/log/auto-translate-bot.log

[Install]
WantedBy=multi-user.target
